# ops/systemd/cajeer-bots@.service
[Unit]
Description=Cajeer Bots — режим %i
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cajeer-bots
EnvironmentFile=/opt/cajeer-bots/.env
ExecStart=/opt/cajeer-bots/.venv/bin/cajeer-bots run %i
Restart=always
RestartSec=5
User=cajeerbots
Group=cajeerbots

[Install]
WantedBy=multi-user.target
