# ops/systemd/cajeer-bots.service
[Unit]
Description=Cajeer Bots — универсальная платформа ботов
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cajeer-bots
EnvironmentFile=/opt/cajeer-bots/.env
ExecStart=/opt/cajeer-bots/.venv/bin/python -m core run all
Restart=always
RestartSec=5
User=cajeerbots
Group=cajeerbots

[Install]
WantedBy=multi-user.target
